# ssh-loop.service
[Unit]
Description=Remote port forward
After=network.target

[Service]
ExecStart=/root/ssh-loop.sh

[Install]
WantedBy=multi-user.target
